# .gemini/skills/agv-mm-arb/configs/supervisord.conf
; ============================================================
; AGV S5 MarketMaker-Agent — supervisord 配置
; D4 决策: 单机双 Campaign 进程管理
; ============================================================
;
; 用法:
;   cd /workspaces/AGV
;   supervisord -c .gemini/skills/agv-mm-arb/configs/supervisord.conf
;
; 管理:
;   supervisorctl -c .gemini/skills/agv-mm-arb/configs/supervisord.conf status
;   supervisorctl -c .gemini/skills/agv-mm-arb/configs/supervisord.conf restart arb-campaign
;   supervisorctl -c .gemini/skills/agv-mm-arb/configs/supervisord.conf stop all
;
; 前提:
;   pip install supervisor          (WQ-YI .venv 已含)
;   .env.s5 已配置                  (BSC RPC + 钱包私钥 + 通知渠道)
;   PYTHONPATH 指向 skills 目录     (下方 environment 已设置)
;

[supervisord]
logfile = /workspaces/AGV/logs/s5/supervisord.log
logfile_maxbytes = 10MB
logfile_backups = 3
loglevel = info
pidfile = /workspaces/AGV/logs/s5/supervisord.pid
nodaemon = false
directory = /workspaces/AGV

; ── Arb Campaign (优先上线 — D2) ──────────────────────────────
[program:arb-campaign]
command = /workspaces/WQ-YI/.venv/bin/python -m _shared.cli.arb_campaign --live --dry-run
directory = /workspaces/AGV
environment =
    PYTHONPATH="/workspaces/AGV/.gemini/skills",
    DOTENV_PATH="/workspaces/AGV/.env.s5"
autostart = true
autorestart = true
startsecs = 10
startretries = 3
stopwaitsecs = 30
redirect_stderr = true
stdout_logfile = /workspaces/AGV/logs/s5/arb-campaign.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 5

; ── MM Campaign (被动护盘 — 不急) ─────────────────────────────
[program:mm-campaign]
command = /workspaces/WQ-YI/.venv/bin/python -m _shared.cli.arb_campaign --campaign=mm --live
directory = /workspaces/AGV
environment =
    PYTHONPATH="/workspaces/AGV/.gemini/skills",
    DOTENV_PATH="/workspaces/AGV/.env.s5"
autostart = false
autorestart = true
startsecs = 10
startretries = 3
stopwaitsecs = 30
redirect_stderr = true
stdout_logfile = /workspaces/AGV/logs/s5/mm-campaign.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 5

; ── supervisorctl 接口 ────────────────────────────────────────
[unix_http_server]
file = /workspaces/AGV/logs/s5/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///workspaces/AGV/logs/s5/supervisor.sock
